# templates/bootcblade-rebuild.timer.j2
[Unit]
Description=bootcblade-rebuild timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
