# templates/bootcblade-rebuild.timer.j2
[Unit]
Description=bootcblade-rebuild timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
